# systemd/usbmuxd.service.in
[Unit]
Description=Socket daemon for the usbmux protocol used by Apple devices

[Service]
Type=simple
ExecStart=@sbindir@/usbmuxd --user usbmux --systemd
PIDFile=@localstatedir@/run/usbmuxd.pid
